# crontab
#run python script every minutes
*/15 13-20 * * * python3 /app/covid-19_crawler/main.py video >> /app/logs/video.log
00 07-24 * * * python3 /app/covid-19_crawler/main.py news >> /app/logs/news.log
